# dir/etc/systemd/system/named.service.d/breakdown.conf
[Service]
ExecStartPre=/mss/share/shart_net-handler named
