# dir/etc/systemd/system/named.service.d/breakdown.conf
[Service]
ExecStartPre=sh -c 'mkdir /tmp/bind/'
ExecStartPre=sh -c 'echo "zone \".\" { type hint; file \"/var/named/zone.root-nov6\"; };" > "/tmp/bind/named.conf.local"'
