# dir/etc/systemd/system/named.service.d/breakdown.conf
[Service]
ExecStartPre=mkdir -p "/tmp/bind"
ExecStartPre=echo "zone \".\" { type hint; file \"/var/named/zone.root-nov6\"; };" > "/tmp/bind/named.conf.local"
