# dir/etc/systemd/system/docker.service.d/netns.conf
[Unit]
JoinsNamespaceOf=dockerns.service
BindsTo=dockerns.service
Requires=dockerns.service
After=dockerns.service

[Service]
PrivateNetwork=Yes
PrivateMounts=No
