# dir/etc/systemd/system/libvirtd.service.d/netns.conf
[Unit]
JoinsNamespaceOf=libvirtdns.service
BindsTo=libvirtdns.service
Requires=libvirtdns.service
After=libvirtdns.service

[Service]
PrivateNetwork=Yes
PrivateMounts=No
